# wb-back/cron.d/laravel-schedule
* * * * * cd /var/www/wb-back && php artisan schedule:run --verbose >> /var/www/wb-back/storage/logs/scheduler.log 2>&1
