# wb-back/cron.d/laravel-schedule
# Laravel Scheduler
* * * * * cd /var/www/wb-back && php artisan schedule:run --verbose >> /var/www/wb-back/storage/logs/scheduler.log 2>&1

# Branches sync - hourly
0 * * * * cd /var/www/wb-back && php artisan app:fetch-branches-info >> /var/www/wb-back/storage/logs/branches-sync.log 2>&1

# Full branches sync - daily at 3 AM
0 3 * * * cd /var/www/wb-back && php artisan app:fetch-branches-info --full >> /var/www/wb-back/storage/logs/branches-sync-daily.log 2>&1
